# ucas-checkin.timer
[Unit]
Description=UCAS Checkin Timer

[Timer]
OnCalendar=*-*-* 6:00
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
