# ucas-checkin.service
[Unit]
Description=UCAS Daily Checkin Service
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/ucas-checkin.py
Environment=PYTHONUNBUFFERED=1
